# backend/crontab
0 */6 * * * root /usr/bin/python /srv/src/manage.py get_balance_commerzbank
0 * * * * root /usr/bin/python /srv/src/manage.py get_balance_n26
